# com.local.capslock-untoggle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.capslock-untoggle</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path after cargo build --release -->
        <string>/usr/local/bin/capslock-untoggle</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Log output in case something goes wrong -->
    <key>StandardOutPath</key>
    <string>/tmp/capslock-untoggle.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/capslock-untoggle.err</string>
</dict>
</plist>
